# rc.local
#!/bin/sh -e
su - pi -c "/home/pi/cookie/bin/python /home/pi/assistant.py >> /home/pi/assistant.log 2>&1 &"
exit 0
